# deploy/vps/maverick.service
# NOTE: the %i / /home/%i tokens below are PLACEHOLDERS. deploy/vps/install.sh
# renders them to the install user + their home before writing this file to
# /etc/systemd/system/. Don't `cp` this file directly -- run install.sh (or
# sed-substitute %i yourself), or systemd gets an empty User= and broken paths.
[Unit]
Description=Maverick agent (long-horizon multi-agent swarm)
Documentation=https://github.com/cdayAI/maverick
After=network-online.target
Wants=network-online.target

# Crash-loop protection: if Maverick crashes 5x in 5 minutes, give up
# instead of respawning forever and burning API credits on a bad config.
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/maverick-workspace
EnvironmentFile=/home/%i/.maverick/.env
Environment=MAVERICK_CONFIG=/home/%i/.maverick/config.toml
# Idle by default — channels (telegram etc.) wake it.
ExecStart=/home/%i/.local/bin/maverick serve
Restart=on-failure
RestartSec=5s
# Give in-flight goals 30s to wind down on stop/restart.
TimeoutStopSec=30
KillMode=mixed

# Security hardening.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/%i/.maverick /home/%i/maverick-workspace
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
